# systemd/home-wordpress-backup.service
[Unit]
Description=Backup home WordPress databases and uploads
Wants=network-online.target
After=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/home-wordpress
ExecStart=/opt/home-wordpress/scripts/backup.sh

