# systemd/home-wordpress-kb-relay-backup.service
[Unit]
Description=Backup relayed personal knowledge base WordPress database and uploads
Wants=network-online.target
After=network-online.target docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/home-wordpress
Environment=KB_COMPOSE_FILE=compose.kb-home-relay.yml
ExecStart=/opt/home-wordpress/scripts/backup-kb.sh
